# config/scripts/bing-sddm-background.timer
[Unit]
Description=Run Bing SDDM Background Update Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
